# deployment/brawlsoccer.service
# Service systemd pour BrawlSoccer
# /etc/systemd/system/brawlsoccer.service

[Unit]
Description=BrawlSoccer - Jeu de Football-Boxe Multijoueur 3D
Documentation=https://github.com/mayeulperrin/brawlsoccer
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/brawlsoccer
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
RestartSec=5
Restart=always

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/brawlsoccer/logs
MemoryLimit=512M
TasksMax=100

# Logs
StandardOutput=append:/var/log/brawlsoccer/app.log
StandardError=append:/var/log/brawlsoccer/error.log
SyslogIdentifier=brawlsoccer

[Install]
WantedBy=multi-user.target